# server_files/marvel.conf
[program:marvel]
command = /opt/.env/bin/gunicorn_start
user = django
stdout_logfile = /opt/logs/marvel_supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

